# conf/rc.local
/data/lampstack/ctlscript.sh  start
supervisord -c /data/supervisord.conf